# debian/syslog-ng-core.conffiles
/etc/systemd/system/multi-user.target.wants/syslog-ng.service
/etc/systemd/system/syslog.service
