# debian/syslog-ng-core.postinst
#! /bin/sh

set -e

if [ "$1" = "triggered" ]; then
   invoke-rc.d syslog-ng stop || exit $?
   invoke-rc.d syslog-ng start || exit $?
   exit 0
fi

dpkg-trigger register-syslog-ng-plugin

#DEBHELPER#

exit 0
