# debian/syslog-ng-core.postinst
#! /bin/sh

set -e

if [ -x "/etc/init.d/syslog-ng" ] && [ ! -n "$2" ]; then
   invoke-rc.d syslog-ng stop || exit $?
fi

#DEBHELPER#

exit 0
